# deploy/supervisor/rost-ai.conf
[program:rost-ai]
directory=/home/rost/new/rost-ai
command=/home/rost/new/rost-ai/scripts/run_prod.sh
user=rost
autostart=true
autorestart=true
startsecs=5
stopsignal=TERM
stdout_logfile=/var/log/rost_ai.log
stderr_logfile=/var/log/rost_ai.err.log
environment=ROST_AI_MODEL_DIR="/home/rost/new/rost-ai/models/best_openvino_model",ROST_AI_ATTACHMENTS_ROOT="/home/rost/files/doc_attachments",ROST_AI_UPLOADS_ROOT="/home/rost/new/rost-ai/runtime/uploads",ROST_AI_HOST="127.0.0.1",ROST_AI_PORT="18080",ROST_AI_WORKERS="1",ROST_AI_INFERENCE_WORKERS="2",ROST_AI_INFERENCE_QUEUE_SIZE="32",PYTHONUNBUFFERED="1"
